# advanced/pihole-update.service
[Unit]
Description=Raspberry Pi Ad Blocker

[Service]
Type=simple
ExecStart=/usr/local/bin/gravity.sh
